# base_project/base_project/supervisor_settings.conf
[program:base_project]
command=$PROJECT_DIR/bin/gunicorn base_project.wsgi:application
directory=$PROJECT_DIR
user=$INSTALL_USER
autostart=true
autorestart=true
redirect_stderr=true
